# config/systemd/tsunami-celery.service
[Unit]
Description=TSUNAMI Celery Worker
After=network.target redis.service
Requires=redis.service

[Service]
Type=simple
User=tsunami
Group=tsunami
WorkingDirectory=/opt/tsunami
Environment="TSUNAMI_ENV=production"
EnvironmentFile=/opt/tsunami/.env

ExecStart=/opt/tsunami/venv/bin/celery -A celery_worker worker \
    --loglevel=info \
    --concurrency=4 \
    --logfile=/var/log/tsunami/celery.log \
    --pidfile=/var/run/tsunami/celery.pid

ExecStop=/bin/kill -s TERM $MAINPID

Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536
MemoryMax=2G

[Install]
WantedBy=multi-user.target
